# packaging/systemd/flightctl-onboarding-mask-greenboot.service
[Unit]
Description=Mask greenboot health checks during system onboarding
Before=greenboot-healthcheck.service
After=local-fs.target
DefaultDependencies=no
ConditionPathExists=!/var/lib/flightctl-onboarding/.onboarding-complete

[Service]
Type=oneshot
ExecStart=/usr/libexec/flightctl-onboarding/mask-greenboot.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
